# .config/systemd/user/fix-virtual-mic.service
[Unit]
Description=Fix Virtual Mic Routing Links
After=pipewire.service wireplumber.service
Requires=pipewire.service

[Service]
Type=oneshot
ExecStart=/home/fib/.config/pipewire/fix-virtual-mic-links.sh
Restart=no

[Install]
WantedBy=pipewire.service
