# .config/systemd/user/pipewire-links.service
[Unit]
Description=PipeWire automatic links and monitoring
After=pipewire.service
Requires=pipewire.service
PartOf=pipewire.service

[Service]
Type=simple
ExecStart=%h/.config/pipewire/links.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=pipewire.service
